# crontab
0 4 * * * sh /update-bt-tracker.sh
# end line
