# crontab
* * * * * sh /update-bt-tracker.sh
# end line
